# deploy/player-web-prod/player-web-backend.service.template
[Unit]
Description=player-web backend
After=network.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__ROOT_DIR__/player-web/server
Environment=PYTHONUNBUFFERED=1
ExecStart=__ROOT_DIR__/player-web/server/.venv/bin/gunicorn --workers 1 --bind 127.0.0.1:8787 --timeout 300 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
